# deb/sysadmws-utils-v1/opt/sysadmws/disk_alert/disk_alert.cron
# Disk Alert
*/5 * * * *	root	/opt/sysadmws-utils/disk_alert/disk_alert.sh
