# mysql_replica_checker/mysql_replica_checker.cron
# Watch MYSQL replication
*/30 * * * *	root	/opt/sysadmws-utils/mysql_replica_checker/mysql_replica_checker.sh
