# mysql_replica_checker/mysql_replica_checker.cron
# Check MYSQL replication
*/30 * * * *	root	/opt/sysadmws/mysql_replica_checker/mysql_replica_checker.sh
