# mysql_replica_checker/replica_checker.cron
# Watch MYSQL replication
*/30 * * * *	root	/opt/sysadmws-utils/watchog/replica/replica_checker.sh >> /opt/sysadmws-utils/watchog/replica/replica_checker.dummy_log 2>&1 
