# deb/sysadmws-utils/etc/cron.d/sysadmws-disk-alert
# Disk Alert
0 10 * * *	root	/opt/sysadmws-utils/disk_alert/disk_alert.sh
