# deb/sysadmws-utils/etc/cron.d/sysadmws-mysql-replica-checker
# Watch MYSQL replication
*/30 * * * *	root	/opt/sysadmws-utils/mysql_replica_checker/mysql_replica_checker.sh | /opt/sysadmws-utils/notify_devilry/notify_devilry.py
