# mysql_queries_log/mysql_queries_log.cron
# Mysql Queries Collection
* * * * *	root	/opt/sysadmws-utils/mysql_queries_log/mysql_queries_log.sh >> /opt/sysadmws-utils/mysql_queries_log/mysql_queries_log.log 2>&1
