# heartbeat_mesh/sender.cron
# heartbeat_mesh/sender
* * * * * /opt/sysadmws/heartbeat_mesh/sender.py
