# bulk_log/bulk_log.cron
# Bulk Log Collection
*/2 * * * *	root	/opt/sysadmws/bulk_log/bulk_log.sh >> /opt/sysadmws/bulk_log/bulk_log.log
