# disk_alert/disk_alert.cron
# Disk Alert
0 10 * * *	root	/opt/sysadmws-utils/disk_alert/disk_alert.sh | /opt/sysadmws-utils/notify_devilry/notify_devilry.py
