# rsnapshot_backup/rsnapshot_backup_daily.cron
30 3  	* * *		root	/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.sh sync >>/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.log 2>&1 && /opt/sysadmws/rsnapshot_backup/rsnapshot_backup.sh daily >>/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.log 2>&1
0  3  	* * 1		root	/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.sh weekly >>/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.log 2>&1
30 2  	1 * *		root	/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.sh monthly >>/opt/sysadmws/rsnapshot_backup/rsnapshot_backup.log 2>&1
