# cmd_check_alert/cmd_check_alert.cron
# cmd_check_alert
*/5 * * * * /opt/sysadmws/cmd_check_alert/cmd_check_alert.sh
